# docker/php7.4-fpm/crontab
* * * * * bash "-c" "/usr/local/bin/php -q -f /var/www/html/bin/console app:news-download" >> /proc/1/fd/1
