# docker/php7.4-fpm/cron.d/crontab
* * * * * bash "-c" "php /var/www/html/bin/console app:news-download" >> /var/www/html/var/log/news_download.log
